# deploy/hostinger/acb-litellm.service
[Unit]
Description=CommandCenter LiteLLM Proxy
After=network-online.target acb.service
Wants=network-online.target
Requires=acb.service

[Service]
Type=simple
User=acb
WorkingDirectory=/opt/acb/app
Environment=PATH=/home/acb/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/opt/acb/app/infra/litellm/litellm.env
ExecStart=/home/acb/.local/bin/uv run litellm --config /opt/acb/app/infra/litellm/config.live.yaml --port 4000
Restart=on-failure
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
